# supervisor/flask-server.conf
[program:flask-server]
command = /home/pi/.virtualenvs/blog/bin/python /home/pi/.virtualenvs/blog/src/flask/motion-socketio/run.py
autostart=true
autorestart=false
user = pi
