# supervisor/motion.conf
[program:motion]
command = motion -d 5 -l /home/pi/motion.log
autostart = true
autorestart = true
