# supervisor/octoprint.conf
[program:octoprint]
command = /home/pi/temp/OctoPrint/venv/bin/octoprint serve
autostart=true
autorestart=false
user = pi
